# pkg/DEBIAN/postinst
#!/bin/sh
set -e

sudo systemctl daemon-reload
systemctl enable cooler || echo "Can't enable service"
systemctl start cooler || echo "Can't start service"
