# pkg-src/DEBIAN/postinst
#!/bin/env sh
set -e

sudo systemctl daemon-reload
systemctl start cooler || echo "Can't start service"
systemctl enable cooler || echo "Can't enable service"
