# pkg-src/DEBIAN/postinst
#!/bin/sh

set -e

pip3 install {{python-packages}}
sudo systemctl daemon-reload
systemctl start raspi-fan-control
systemctl enable raspi-fan-control
