# pkg-src/DEBIAN/postrm
#!/bin/sh

set -e

systemctl disable raspi-fan-control
systemctl stop raspi-fan-control
rm /usr/bin/raspi-fan-control
