# pkg-src/DEBIAN/prerm
#!/bin/sh

set -e

systemctl disable raspi-fan-control || echo "Can't disable service"
systemctl stop raspi-fan-control || echo "Can't stop service"
