# pkg-src/DEBIAN/prerm
#!/bin/env sh
set -e

systemctl disable cooler || echo "Can't disable service"
systemctl stop cooler || echo "Can't stop service"
